# configs/etc/supervisor/conf.d/geoad.conf
[program:geoad]
command=npm run server
directory=/var/www/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/geoad.err.log
stdout_logfile=/var/log/geoad.out.log
user=root
